# Backend/ServerForTwitterStream/scrape.cron
*/15 * * * * python3 TwitterScraperS3_electionnight.py
